# python/tg-daemon-watcher.path
[Path]
Unit=tg-daemon-watcher.service
PathChanged=/home/bitrix/ext_www/domain.tld/daemon
# trigger on changes to a file not just create/delete
# don't put comments in the same line as the command

[Install]
WantedBy=multi-user.target
